# dotfiles/smart_plug/linux/speakers.service
[Unit]
Description=Speakers
Wants=network-online.target
After=network-online.target
Conflicts=sleep.target
Before=sleep.target

[Service]
Type=oneshot
ExecStart=/opt/speaker/speakerctl on
RemainAfterExit=yes
ExecStop=/opt/speaker/speakerctl off

[Install]
WantedBy=multi-user.target wake.target
