# program-configs/nix/non-nixos-gpu-setup/non-nixos-gpu-x.service
[Unit]
Description=Set up GPU libs

[Service]
Type=oneshot
ExecStart=@setupbash@
