# program-configs/nix/non-nixos-gpu-setup/non-nixos-gpu-x.path
[Unit]
Description=Trigger non-nixos GPU setup whenever dpkg does something

[Path]
PathChanged=/var/log/dpkg.log

[Install]
WantedBy=multi-user.target
