# program-configs/smart-plug/linux/wake.target
[Unit]
Description=Wake
Conflicts=sleep.target
Before=sleep.target
